# devops_config/supervisor.ini
[program:app]
directory=/code
command = /usr/local/bin/uwsgi --ini /code/backend_uwsgi.ini
autostart = true
autorestart = true
exitcodes = 0,2
stopsignal = TERM
redirect_stderr = true
stdout_logfile = /var/log/app.log

[program:celery]
directory=/code
command = /usr/local/bin/celery worker -A GolfConnect -B --loglevel=ERROR -s /code/media/celerybeat-schedule
autostart = true
autorestart = true
exitcodes = 0,2
stopsignal = TERM
redirect_stderr = true
stdout_logfile = /var/log/celery.log
